# infra/systemd/xignal-backend.service.example
[Unit]
Description=Xignal Backend
After=network.target

[Service]
WorkingDirectory=/srv/xignal
EnvironmentFile=-/etc/xignal/backend.env
ExecStart=/usr/bin/node legacy/backend/bin/www
Environment=PORT=3000
Environment=NODE_ENV=production
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
StandardOutput=append:/var/log/xignal/backend-jsonl.log
StandardError=append:/var/log/xignal/backend-jsonl.log

[Install]
WantedBy=multi-user.target
